# daemon/bright.service
[Unit]
Description=Bright
After=network.target

[Service]
Type=simple
# UPDATE THIS: Your username
User=YOUR_USERNAME
# UPDATE THIS: Path to relay directory
WorkingDirectory=/home/YOUR_USERNAME/bright
# UPDATE THIS: Path to bun
ExecStart=/home/YOUR_USERNAME/.bun/bin/bun run src/relay.ts
Restart=always
RestartSec=10

# Environment
Environment=HOME=/home/YOUR_USERNAME
Environment=PATH=/home/YOUR_USERNAME/.bun/bin:/usr/local/bin:/usr/bin:/bin

# Optional: Load .env file
# EnvironmentFile=/home/YOUR_USERNAME/bright/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bright

[Install]
WantedBy=multi-user.target

# SETUP INSTRUCTIONS:
#
# 1. Replace all YOUR_USERNAME with your actual Linux username
#
# 2. Copy to systemd:
#    sudo cp daemon/bright.service /etc/systemd/system/
#
# 3. Reload systemd:
#    sudo systemctl daemon-reload
#
# 4. Enable (start on boot):
#    sudo systemctl enable bright
#
# 5. Start:
#    sudo systemctl start bright
#
# 6. Check status:
#    sudo systemctl status bright
#
# 7. View logs:
#    journalctl -u bright -f
#
# 8. Stop:
#    sudo systemctl stop bright
#
# 9. Disable:
#    sudo systemctl disable bright
